# debian/postinst
#!/bin/sh
set -e

case "$1" in
    configure)
        systemctl daemon-reload
        systemctl enable pifan
        systemctl start pifan
        ;;
esac
